# postgis94.service
[Unit]
Description=Postgis94
After=docker.service
Requires=docker.service

[Service]
ExecStart=/opt/bin/systemd-docker run --rm --name %n -v /btrfs/postgis94:/var/lib/postgresql/data -e POSTGRES_PASSWORD="Aev6ieL8aequekue" mdillon/postgis:9.4
Restart=always
RestartSec=10s
Type=notify
NotifyAccess=all
TimeoutStartSec=120
TimeoutStopSec=15

[Install]
WantedBy=multi-user.target
